# cmder/config/profile.d/directory-update.sh
#!/usr/bin/env bash
#
# Notify conemu when we switch directories (updates the tab title appropriately)
PROMPT_COMMAND='ConEmuC -StoreCWD'
